# VOLUMIO-LCD/systemd/volumio-lcd.service
[Unit]
Description=Volumio LCD album art display
After=network-online.target volumio.service
Wants=network-online.target

[Service]
Type=simple
User=volumio
WorkingDirectory=/home/volumio/waveshare-2.8/Python
ExecStart=/usr/bin/python3 /home/volumio/waveshare-2.8/Python/volumio_lcd.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
