# VOLUMIO-LCD/systemd/volumio-lcd.service
[Unit]
Description=Volumio LCD display service
After=network.target
Wants=network.target

[Service]
Type=simple
User=volumio
Group=volumio
WorkingDirectory=/home/volumio/waveshare-2.8/Python
ExecStart=/usr/bin/python3 /home/volumio/waveshare-2.8/Python/volumio_lcd.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

